# frauzufall-walls.conf
 #!upstart
description "node.js server for walls, an interactive projection by frauzufall"

start on started mountall
stop on shutdown

#Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

script
        cd /root/apps/crawling-traces-server/
        node app.js >> logs/walls.out
end script
